# telegram-monitor.service
[Unit]
Description=Telegram Channel Monitor Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/test-telegram-bot
Environment="PATH=/root/test-telegram-bot/venv/bin:/usr/local/bin:/usr/bin:/bin"
# Явно указываем путь к .env файлу
Environment="PYTHONUNBUFFERED=1"
ExecStart=/root/test-telegram-bot/venv/bin/python /root/test-telegram-bot/main.py
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=10
# Увеличиваем время ожидания перед перезапуском
StartLimitInterval=600
StartLimitBurst=10
# Перезапускаем при любых ошибках
RestartPreventExitStatus=

[Install]
WantedBy=multi-user.target
